# mjpg-streamer-experimental/debian/mjpg-streamer@.service
[Unit]
Description=A server for streaming Motion-JPEG from a video capture device
After=network.target
Wants=dev-%I.device
After=dev-%I.device

[Service]
User=mjpg_streamer
ExecStart=/usr/bin/mjpg_streamer -i 'input_uvc.so -d /dev/%I' -o 'output_http.so -w /usr/share/mjpg-streamer/www'

[Install]
WantedBy=multi-user.target
